# deploy/trader3-bot.service
[Unit]
Description=Trader3 Bot (paper trading Testnet)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/app-trader
# Dùng Python trong môi trường ảo .venv
ExecStart=/home/ubuntu/app-trader/.venv/bin/python run_paper.py
Restart=always
RestartSec=15
# Ghi output ra file để xem lại
StandardOutput=append:/home/ubuntu/app-trader/service.log
StandardError=append:/home/ubuntu/app-trader/service.log

[Install]
WantedBy=multi-user.target
